# homelab-screen.service
[Unit]
Description=Thermalright AIO Cooler USB LCD System Monitor
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/homelab-screen
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=false
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
SystemCallArchitectures=native
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
